# projects/t31_lite_convoy/general/scripts/excludes/t31_lite.list
etc/sensor/gc2083-t31.bin
etc/sensor/gc2083.yaml
lib/modules/3.10.14__isvp_swan_1.0__/extra/sensor_gc2083_t31.ko
etc/sensor/gc4023-t31.bin
etc/sensor/gc4023.yaml
lib/modules/3.10.14__isvp_swan_1.0__/extra/sensor_gc4023_t31.ko
etc/sensor/gc4653-t31.bin
etc/sensor/gc4653.yaml
lib/modules/3.10.14__isvp_swan_1.0__/extra/sensor_gc4653_t31.ko
etc/sensor/imx327-t31.bin
etc/sensor/imx327.yaml
lib/modules/3.10.14__isvp_swan_1.0__/extra/sensor_imx327_t31.ko
etc/sensor/jxf37-t31.bin
etc/sensor/jxf37.yaml
lib/modules/3.10.14__isvp_swan_1.0__/extra/sensor_jxf37_t31.ko
etc/sensor/jxh62-t31.bin
etc/sensor/jxh62.yaml
lib/modules/3.10.14__isvp_swan_1.0__/extra/sensor_jxh62_t31.ko
etc/sensor/jxh63.yaml
lib/modules/3.10.14__isvp_swan_1.0__/extra/sensor_jxh63_t31.ko
etc/sensor/jxq03p-t31.bin
etc/sensor/jxq03p.yaml
lib/modules/3.10.14__isvp_swan_1.0__/extra/sensor_jxq03p_t31.ko
etc/sensor/jxq03-t31.bin
etc/sensor/jxq03.yaml
lib/modules/3.10.14__isvp_swan_1.0__/extra/sensor_jxq03_t31.ko
etc/sensor/sc200ai-t31.bin
etc/sensor/sc200ai.yaml
lib/modules/3.10.14__isvp_swan_1.0__/extra/sensor_sc200ai_t31.ko
etc/sensor/sc2232h-t31.bin
etc/sensor/sc2232h.yaml
lib/modules/3.10.14__isvp_swan_1.0__/extra/sensor_sc2232h_t31.ko
etc/sensor/sc2335-t31.bin
etc/sensor/sc2335.yaml
lib/modules/3.10.14__isvp_swan_1.0__/extra/sensor_sc2335_t31.ko
etc/sensor/sc2336-t31.bin
etc/sensor/sc2336.yaml
lib/modules/3.10.14__isvp_swan_1.0__/extra/sensor_sc2336_t31.ko
etc/sensor/sc3335-t31.bin
etc/sensor/sc3335.yaml
etc/sensor/sc3338-t31.bin
etc/sensor/sc3338.yaml
lib/modules/3.10.14__isvp_swan_1.0__/extra/sensor_sc3338_t31.ko
etc/sensor/sc4236-t31.bin
etc/sensor/sc4236.yaml
lib/modules/3.10.14__isvp_swan_1.0__/extra/sensor_sc4236_t31.ko
etc/sensor/sc5235-t31.bin
etc/sensor/sc5235.yaml
lib/modules/3.10.14__isvp_swan_1.0__/extra/sensor_sc5235_t31.ko
etc/sensor/jxf22.yaml
etc/sensor/jxf23.yaml
etc/sensor/jxh42.yaml
etc/sensor/jxk04.yaml
etc/sensor/os03b10.yaml
etc/sensor/ov2735b.yaml
etc/sensor/ov2735.yaml
etc/sensor/ov4689.yaml
etc/sensor/ps5260.yaml
etc/sensor/sc2232.yaml
etc/init.d/S49ntpd
etc/init.d/S50snmpd
etc/init.d/S50telnet
etc/init.d/S60crond
etc/init.d/S92motion
etc/init.d/S93telegrambot
etc/init.d/S94rc.local